# tarsnap/.config/systemd/user/tarsnap.service
[Unit]
Description=Run tarsnap backup

[Service]
Type=oneshot
ExecStart=/bin/sh ${HOME}/.config/tarsnap/backup.sh
