# backup/.config/systemd/user/backup.service
[Unit]
Description=Run backup

[Service]
Type=oneshot
ExecStart=/bin/fish backup-run
