# backup/.config/systemd/user/backup.service
[Unit]
Description=Run backup

[Service]
Type=oneshot
ExecStart=/usr/bin/fish -c backup-run
