# scrapers/.config/systemd/user/scrapers.service
[Unit]
Description=Scrape sites

[Service]
Type=oneshot
WorkingDirectory=/home/dave/Workspace/davebrent/scrapers
ExecStart=make scrape
